# openhasp-fbdev/openhasp-fbdev.initd
#!/sbin/openrc-run

name="openHASP fbdev"
command="/usr/bin/openhasp-fbdev"
command_args="-c /usr/share/hasp"
supervisor="supervise-daemon"

